# docker/crontab
# 毎時0分に実行
0 * * * * cd /app/auto_post_bot && python post_tweet.py >> /var/log/cron.log 2>&1

